# .config/systemd/user/my-graphical-session.target
[Unit]
Description=My Graphical Session
Requires=graphical-session-pre.target
BindsTo=graphical-session.target
