# .config/systemd/user/xfsettingsd.service
[Unit]
Description=XFCE settings daemon
PartOf=my-graphical-session.target

[Service]
ExecStart=/usr/bin/env xfsettingsd

[Install]
WantedBy=my-graphical-session.target
